# home/infrastructure_as_code/ansible/taskset_ansible_infrastructure_as_code/task_008_kibana_docker/playbooks/roles/drone-ci/templates/system/drone.service
[Unit]
Description=Drone server
After=docker.service nginx.service

[Service]
Restart=always
ExecStart=/usr/local/bin/docker-compose -f /etc/drone/docker-compose.yml up
ExecStop=/usr/local/bin/docker-compose -f /etc/drone/docker-compose.yml stop

[Install]
WantedBy=multi-user.target